# etc/skel/.config/autostart/nitroshare.desktop
[Desktop Entry]
Version=1.0
Name=NitroShare
Type=Application
Exec=nitroshare
Terminal=false
